# docker/holyclaude/scripts/s6/claude-mem-worker/run
#!/command/with-contenv sh
# ==============================================================================
# s6-overlay longrun: claude-mem worker (HTTP API + observation capture).
#
# Supervised by s6 — auto-restarted on crash, started at container boot.
# The entrypoint wrapper provisions the plugin (marketplace + enabledPlugins)
# BEFORE s6 starts, so the worker finds claude-mem already installed.
#
# Runtime config is supplied via environment (docker-compose):
#   CLAUDE_MEM_MODEL        summarization model (glm-5.2 — Z.AI routing)
#   CLAUDE_MEM_WORKER_HOST  0.0.0.0 so the port-forward reaches the host browser
#   CLAUDE_MEM_WORKER_PORT  pinned to 37700 (default drifts with uid)
# ==============================================================================
export HOME=/home/claude
CMEM_VERSION="${CLAUDE_MEM_VERSION:-13.9.1}"
export CMEM_VERSION

worker_running() {
  ps -eo args | grep -F "worker-service.cjs --daemon" | grep -v grep >/dev/null 2>&1
}

# claude-mem v13 starts the actual worker as a daemonized plugin process. Keep
# this s6 longrun alive as a small monitor so s6 does not restart-loop npx while
# the daemon is already running.
while true; do
  if ! worker_running; then
    s6-setuidgid claude npx -y "claude-mem@${CMEM_VERSION}" start --daemon || true
  fi
  sleep 30
done
